# conf/user/afm-user-start.service
[Unit]
Description=start mark of afm user

[Service]
Type=oneshot
ExecStart=/bin/touch %t/afm-user-started


