# conf/system/afm-user-session@.target
[Unit]
Description=afm-user-session %i
After=user@%i.service
