# conf/system/afm-user-session@.target
[Unit]
After=afm-user-session@%i.service
After=afm-user-setup@%i.service
