# conf/system/afm-user-setup@.service.in
[Unit]
Description=Initiate afm-user-setup %i
After=user-runtime-dir@%i.service
Before=user@%i.service
[Service]
Type=oneshot
StartLimitInterval=0
ExecStart=-@afm_libexecdir@/afm-user-setup.sh %i
