# conf/system/afm-api-afm-main@.service.in
[Unit]
Description=Provides api afm-main for user %i
X-AFM-API-TYPE=ws

Requires=afm-user-setup@%i.service
After=afm-user-setup@%i.service
Requires=afm-system-daemon.socket
After=afm-system-daemon.socket

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/ln -s @afm_platform_rundir@/apis/ws/afm-main @afm_users_rundir@/%i/apis/ws/afm-main
