# conf/system/afm-system-daemon.socket.in

[Socket]
SmackLabel=*
ListenStream=/run/afm/apis/ws/afm-main
FileDescriptorName=afm-main

Service=afm-system-daemon.service

[Install]
WantedBy=sockets.target

