# conf/system/afm-user-session@.service.in
[Unit]
Description=Initiate afm-user-session %i
Requires=systemd-user-sessions.service
After=systemd-user-sessions.service
[Service]
User=%i
PAMName=afm-user-session
ExecStart=@afm_libexecdir@/afm-user-session
AmbientCapabilities=CAP_MAC_OVERRIDE
SecureBits=no-setuid-fixup
