# conf/system/afm-user-session@.target
[Unit]
Description=Target group of afm-user-session %i
Requires=afm-user-setup@%i.service
After=afm-user-setup@%i.service
