# conf/system/afm-system-setup.service.in
[Unit]
Description=Initiate afm-system-setup
Before=weston@display.service afs-supervisor.service afm-system-daemon.service
[Service]
Type=oneshot
ExecStart=-@afm_libexecdir@/afm-system-setup.sh
