# conf/system/afmpkg-installer.service.in
[Unit]
Description=Installer for application framework

[Service]
Type=simple
ExecStart=/usr/bin/afmpkg-installerd --socket sd:sock
Sockets=afmpkg-installer.socket

[Install]
WantedBy=multi-user.target
