# conf/system/afmpkg-installer.service.in
[Unit]
Description=Installer for application framework

[Service]
Type=simple
ExecStart=/usr/bin/afmpkg-installer --socket sd:sock
Sockets=afmpkg-installer.socket
