# tests/test3/ref/ROOT/system/afm-api-toto@.service
# auto generated by redpesk framework for example3 version 1.0 target main of example3
[Unit]
Description=Provides api toto for user %i
X-AFM-API-TYPE=auto
BindsTo=afm-user-setup@%i.service
After=afm-user-setup@%i.service
Requires=afm-api-toto@%i.socket
After=afm-api-toto@%i.socket
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
