# tests/test3/ref/ROOT/system/afm-appli-example3--main@.service
# auto generated by redpesk framework for example3 version 1.0 target main of example3

[Unit]
Description=tiny example\n
X-AFM-description=tiny example\n
X-AFM-name=
X-AFM-shortname=
X-AFM-id=example3
X-AFM-version=1.0
X-AFM-author=
X-AFM-author-email=
X-AFM-width=
X-AFM-height=
X-AFM--ID=11111
X-AFM--redpak-id=[REDPAKID]
X-AFM--rootdir=[ROOTDIR]
X-AFM--wgtdir=[ROOTDIR]/[INSTALLDIR]
X-AFM--workdir=[ROOTDIR]/home/%i/app-data/example3
X-AFM--target-name=main
X-AFM--content=/usr/share/example3/a/f2
X-AFM--type=application/x-executable
X-AFM--visibility=visible

X-AFM--scope=user
BindsTo=afm-user-session@%i.target
After=user@%i.service
After=Network.target
# Adds check to smack or selinux
ConditionSecurity=|smack
ConditionSecurity=|selinux

# Automatic bound to required api
Requires=afm-api-toto@%i.socket
After=afm-api-toto@%i.socket

[Service]
EnvironmentFile=-[ROOTDIR]/CONFDIR/afm/unit.env.d/*
EnvironmentFile=-[ROOTDIR]/CONFDIR/afm/widget.env.d/example3/*
SmackProcessLabel=App:example3
SELinuxContext=system_u:system_r:example3_t:s0
SuccessExitStatus=0 SIGKILL
UMask=0077
User=%i
Slice=user-%i.slice
WorkingDirectory=-[ROOTDIR]/home/%i/app-data/example3
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%i/bus
CapabilityBoundingSet=
ExecStartPre=-/usr/bin/redwrap --redpath [ROOTDIR] --  /bin/mkdir -p /home/%i/app-data/example3
SystemCallFilter=~@clock

Environment=AFM_ID=example3
Environment=AFM_APP_INSTALL_DIR=[INSTALLDIR]
Environment=PATH=/usr/sbin:/usr/bin:/sbin:/bin:[INSTALLDIR]/bin
Environment=LD_LIBRARY_PATH=[INSTALLDIR]/lib
Environment=AFM_WORKDIR=/home/%i/app-data/example3
Environment=AFM_WSAPI_DIR=/run/user/%i/apis/ws
Environment=XDG_DATA_HOME=/home/%i/app-data/example3
Environment=XDG_CONFIG_HOME=/home/%i/app-data/example3
Environment=XDG_CACHE_HOME=/home/%i/app-data/example3
Environment=XDG_RUNTIME_DIR=/run/user/%i
SyslogIdentifier=afbd-example3
StandardInput=null
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/redwrap --redpath [ROOTDIR] --  [INSTALLDIR]/usr/share/example3/a/f2
