# tests/test3/ref/ROOT/system/afm-api-toto@.socket
# auto generated by redpesk framework for example3 version 1.0 target main of example3
[Unit]
Description=Provides websocket api toto for user %i
BindsTo=afm-user-setup@%i.service
After=afm-user-setup@%i.service
DefaultDependencies=no
[Socket]
SmackLabel=*
SmackLabelIPIn=System
SmackLabelIPOut=System
SELinuxContextFromNet=true
ListenStream=/run/user/%i/apis/ws/toto
FileDescriptorName=toto
Service=afm-appli-example3--main@%i.service
